# deploy/systemd/hannario-bot.service.example
[Unit]
Description=Hannario Discord bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/keitaito/hannario-v2
EnvironmentFile=/home/keitaito/hannario-v2/.env
ExecStart=/home/keitaito/.local/bin/uv run python bot.py
Restart=on-failure
RestartSec=10
TimeoutStopSec=30

# Keep Python output visible in journalctl.
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
